# trinity-base/lisa/files/lisa
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
    use net
}

start() {
    ebegin "Starting lisa"
    start-stop-daemon --start -x _TDEDIR_/bin/lisa -q -- $LISACONF 1>/dev/null
    eend $? "Failed to start lisa"
}

stop() {
    ebegin "Stopping lisa"
    start-stop-daemon --stop -x _TDEDIR_/bin/lisa -q
    eend $? "Failed to stop lisa"
}
